# scripts/com.atlas.voice-memo-processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atlas.voice-memo-processor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eriklaine/atlas-dashboard-v2/scripts/process-voice-memos.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/eriklaine/.local/bin</string>
        <key>ATLAS_DASHBOARD_URL</key>
        <string>https://atlas-dashboard-v2-ten.vercel.app</string>
    </dict>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/atlas-voice-processor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/atlas-voice-processor.err</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
